# docker/launchd-agents/com.iap.docker-events.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.iap.docker-events</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/opt/homebrew/bin/docker events --format '{{json .}}' >> /Users/admin/.platform-logs/docker-events.log 2>> /Users/admin/.platform-logs/docker-events.err</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/admin/.platform-logs/docker-events-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/admin/.platform-logs/docker-events-launchd.err</string>
</dict>
</plist>
