# docker/launchd-agents/com.iap.vault-audit-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.iap.vault-audit-rotate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/admin/.platform-scripts/vault-audit-rotate.sh; rc=$?; touch /Users/admin/.platform-logs/vault-audit-rotate.heartbeat; exit $rc</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/admin/.platform-logs/vault-audit-rotate.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/admin/.platform-logs/vault-audit-rotate.err</string>
</dict>
</plist>
