# docker/launchd-agents/com.iap.arr-apikey-sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.iap.arr-apikey-sweep</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/admin/repos/integrated-ai-platform/scripts/arr-apikey-sweep.sh &gt;&gt; /Users/admin/.platform-logs/arr-apikey-sweep.log 2&gt;&amp;1; rc=$?; touch /Users/admin/.platform-logs/arr-apikey-sweep.heartbeat; exit $rc</string>
    </array>
    <!-- Hourly sweep: functional-probe-based Prowlarr application resync (D-17-50, §18.E). -->
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/admin/.platform-logs/arr-apikey-sweep.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/admin/.platform-logs/arr-apikey-sweep.launchd.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
