# deployment/launchd/com.adriancox.litellm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.adriancox.litellm</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/adriancox/.local/bin/litellm</string>
		<string>--config</string>
		<string>/Users/adriancox/local-ai-workstation/configs/litellm/config.yaml</string>
		<string>--port</string>
		<string>4000</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/adriancox/local-ai-workstation/logs/litellm.out</string>
	<key>StandardErrorPath</key>
	<string>/Users/adriancox/local-ai-workstation/logs/litellm.err</string>
	<key>WorkingDirectory</key>
	<string>/Users/adriancox</string>
</dict>
</plist>
